# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/var/log/supervisord.log

[program:data_ingestion]
command=python main.py
directory=/app
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/data_ingestion.out.log
stderr_logfile=/var/log/data_ingestion.err.log

[program:coinbase_retention]
command=/usr/local/bin/supercronic /etc/cron.d/coinbase_retention.cron
autostart=true
autorestart=true
stdout_logfile=/var/log/coinbase_retention.out.log
stderr_logfile=/var/log/coinbase_retention.err.log
