# cron/coinbase_retention.cron
# Runs 15 minutes past every hour.
15 * * * * python storage/retention_cleaner.py --exchange coinbase >> /var/log/coinbase_retention.log 2>&1
