# cron/gap_audit.cron
30 2 * * * python /app/services/gap_audit.py >> /var/log/gap_audit.log 2>&1
