# cron/gen_sync.cron
*/5 * * * * python /app/services/gen_candle_sync.py >> /var/log/gen_sync.log 2>&1
